# deploy/transcript-server.service
# systemd unit for the transcript API server (Linux GPU host).
#
# Install:
#   sudo cp deploy/transcript-server.service /etc/systemd/system/
#   sudo install -m 600 deploy/transcript-server.env.example /etc/transcript-server.env
#   sudoedit /etc/transcript-server.env
#   sudo systemctl edit transcript-server   # set User and WorkingDirectory below
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now transcript-server
#   journalctl -u transcript-server -f       # logs
#
# Edit the placeholders (USER, /path/to/...) to match your machine.

[Unit]
Description=transcript API server (WhisperX transcription + diarization)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
UMask=0077
User=USER
WorkingDirectory=/path/to/transcript
# Point at the venv that has the [server] extra + torch installed.
Environment=PATH=/path/to/transcript/.venv/bin:/usr/local/bin:/usr/bin:/bin
# Root-owned secrets and runtime configuration.
EnvironmentFile=/etc/transcript-server.env
ExecStart=/path/to/transcript/.venv/bin/transcript-server --host ${TRANSCRIPT_HOST} --port ${TRANSCRIPT_PORT} --model ${TRANSCRIPT_MODEL}
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
